# .config/systemd/user/terminal.service
[Unit]
Description=Terminal
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/autostart-tagged ghostty

[Install]
WantedBy=graphical-session.target